# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]

logfile=/home/griz/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=griz                                       ; default user
childlogdir=/home/griz/supervisord              ; where child log files will live

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:touchscreen]
command=/usr/sbin/touchscreen start-foreground
autostart=true
stopsignal=TERM
startsecs=3
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/touchscreen-stderr.log
stdout_logfile = /home/griz/supervisord/log/touchscreen-stdout.log

[program:rover_roscore]
command=/home/griz/supervisord/startup/rosl_roscore start-foreground
autostart=true
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_roscore-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover-roscore-stdout.log

[program:rover_base]
command=/home/griz/supervisord/startup/rosl_base -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_base-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_base-stdout.log

[program:rover_imu]
command=/home/griz/supervisord/startup/rosl_imu -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_imu-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_imu-stdout.log

[program:rover_tethered]
command=/home/griz/supervisord/startup/rosl_tethered -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_tethered-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_tethered-stdout.log

[program:rover_livox]
command=/home/griz/supervisord/startup/rosl_livox -start_foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_livox-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_livox-stdout.log

[program:rover_livox_nav]
command=/home/griz/supervisord/startup/rosl_livox_nav -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_livox_nav-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_livox_nav-stdout.log

[program:rover_localization]
command=/home/griz/supervisord/startup/rosl_localization -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_localization-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_localization-stdout.log

[program:rover_octomap]
command=/home/griz/supervisord/startup/rosl_octomap -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_octomap-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_octomap-stdout.log

[program:rover_zed]
command=/home/griz/supervisord/startup/rosl_zed -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_zed-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_zed-stdout.log

[program:rover_zed_nav]
command=/home/griz/supervisord/startup/rosl_tethered -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_zed_nav-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_zed_nav-stdout.log

[program:rover_web]
command=/home/griz/supervisord/startup/rosl_web -start-foreground
autostart=false
stopsignal=TERM
stopasgroup=true
stderr_logfile = /home/griz/supervisord/log/rover_web-stderr.log
stdout_logfile = /home/griz/supervisord/log/rover_web-stdout.log
